# desktop/com.posture.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.posture.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Posture Monitor.app/Contents/MacOS/Posture Monitor</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>LaunchOnlyOnce</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>/tmp/posture-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/posture-monitor.error.log</string>
</dict>
</plist>
